# bin/ap_scanner.service
[Unit]
Description=AP Scanner Daemon Service.
After=network.target

[Service]
ExecStart=/usr/bin/sudo /usr/local/bin/ap_scanner -d -p "daemon"
Restart=always
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target